# manifest/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
nginx -g 'pid /var/lib/nginx/nginx.pid; daemon off;'
