# manifest/etc/services.d/raumserver/run
#!/usr/bin/execlineb -P
sh -c "/bin/cp -a /config . && /bin/mkdir logs && /usr/bin/node /node_modules/node-raumserver/raumserver.js"
